# sv/apache/run
#!/bin/sh
exec 2>&1
exec /etc/init.d/apache2 restart
